# service/server.service
[Unit] 
Description=Server Running Backend for Matthew Wen's mwenclubhouse
[Service]
Type = simple
ExecStart = /home/mwen/GitLab/server/server
